# templates/default/upstart.conf.erb
# nomad - nomad job file

description "Nomad System Scheduler"
author "Nathan Williams <nath.e.will@gmail.com>"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart if crashed
respawn

# Start the process
<% if @server %>
exec /usr/local/bin/nomad agent <%= @daemon_args %> -server
<% else %>
exec /usr/local/bin/nomad agent <%= @daemon_args %>
<% end %>
